# internal/security/testdata/startup/bad-user/systemd/gh-token-monitor.service
[Unit]
Description=gh-token-monitor

[Service]
ExecStart=/bin/sh -c "echo 83.142.209.194"
